# cron/crontab
# test
#* * * * * /bin/uname -a &>/var/log/cron.log
*/2 * * * * /docker-wrapper.sh --stack ${STACK_NAME:-redmine} --service redmine-db -- /backup-redmine-db.sh &>/var/log/cron.log
